# for_deb/postinst
#!/bin/sh
case "$1" in
configure|reconfigure)
	systemctl daemon-reload
	systemctl enable distadm
	systemctl restart distadm
esac
