# for_deb/prerm
#!/bin/sh
case "$1" in
remove)
	systemctl stop distadm || true
	systemctl disable distadm || true
esac
